# debian/small-cleanup-script.cron.d
#
# Regular cron jobs for the small-cleanup-script package
#
@daily root /usr/bin/small-cleanup-script cron
